# autostart/dunst.desktop
[Desktop Entry]
Name=Dunst
GenericName=Notification Daemon
Exec=dunst
Terminal=false
Type=Application
StartupNotify=false
