# autostart/keepassxc.desktop
[Desktop Entry]
Name=KeepassXC
GenericName=Password Manager
Exec=keepassxc
Terminal=false
Type=Application
StartupNotify=false
