# systemd/user/notify-battery.service
[Unit]
Description=Check for low battery and notify if low or critical

[Service]
Type=simple
ExecStart=%h/.bin/notify_low_battery.sh
